# internet-gateway/etc/rc.local
#add ipset restore
ipset restore < /etc/chnroute.ipset
